# supervisor.conf.d/dnscrypt.conf
[program:dnscrypt]
command=/usr/local/sbin/dnscrypt-proxy /usr/local/etc/dnscrypt-proxy.conf
user=root
