# supervisor.conf.d/unbound.conf
[program:unbound]
command=/usr/local/sbin/unbound -c /usr/local/etc/unbound/unbound.conf
