# supervisor/conf.d/ubuntu-dnscrypt.conf
[program:unbound]
command=/usr/local/sbin/unbound -c /usr/local/etc/unbound/unbound.conf


[program:dnscrypt]
command=/usr/local/sbin/dnscrypt-proxy /usr/local/etc/dnscrypt-proxy.conf
user=root